# runit/proxngin/run
#!/bin/bash
exec /usr/sbin/proxngin -t /data/template -c /etc/nginx/dynamic.d
